# scripts/launchd-feishu-daily.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  1. 将本文件复制到 ~/Library/LaunchAgents/com.claude-test.feishu-daily.plist
  2. 把下面 ProgramArguments 里的路径改成你电脑上的真实路径
  3. 修改 StartCalendarInterval 的 Hour / Minute（24 小时制，本地时区）
  4. Webhook 写入 ~/.feishu_daily_webhook（仅一行 URL），或改用 EnvironmentVariables（见飞书说明）
  5. launchctl load ~/Library/LaunchAgents/com.claude-test.feishu-daily.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-test.feishu-daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/heyuxian/Desktop/claude-test/scripts/run_feishu_daily.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/heyuxian/Library/Logs/feishu-daily.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/heyuxian/Library/Logs/feishu-daily.err</string>
</dict>
</plist>
